# scripts/outfitpi.service
[Unit]
Description=OutfitPi — kid-friendly weather outfit recommender
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/outfitpi
ExecStart=%h/outfitpi/venv/bin/python %h/outfitpi/app.py
Restart=on-failure
RestartSec=3
SuccessExitStatus=3
Environment=OUTFITPI_LOG_LEVEL=INFO
# `User=` is set by install.sh via sed substitution.
User=__OUTFITPI_USER__

[Install]
WantedBy=default.target
