# scripts/outfitpi.service
[Unit]
Description=OutfitPi — kid-friendly weather outfit recommender
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/outfitpi
ExecStart=%h/outfitpi/venv/bin/python %h/outfitpi/app.py
Restart=always
RestartSec=3
SuccessExitStatus=3
Environment=OUTFITPI_LOG_LEVEL=INFO

[Install]
WantedBy=default.target
